# scripts/systemd/rangebar-kintsugi.service
# Issue #115, #117: Kintsugi self-healing gap reconciler with watchdog
[Unit]
Description=RangeBar Kintsugi (Self-Healing Gap Reconciler)
After=network-online.target clickhouse-server.service
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
# Issue #117: sd_notify("WATCHDOG=1") sent after each kintsugi_pass()
Type=notify
NotifyAccess=all
WatchdogSec=3600
WorkingDirectory=/home/tca/rangebar-py
EnvironmentFile=/home/tca/rangebar-py/deploy/rangebar.env
ExecStart=/home/tca/rangebar-py/.venv/bin/python3 /home/tca/rangebar-py/scripts/kintsugi.py --daemon

# Restart policy
Restart=always
RestartSec=60

# Resource limits (lower priority than sidecar)
Nice=5
MemoryHigh=3G
MemoryMax=4G
MemorySwapMax=0
OOMScoreAdjust=-200

# Process limits
LimitNOFILE=65536
LimitNPROC=4096

# Standard output/error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rangebar-kintsugi

[Install]
WantedBy=default.target
