# scripts/systemd/rangebar-sidecar.service
[Unit]
Description=rangebar-py streaming sidecar (Issue #107, #117)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
# Issue #117: healthdog-rs bridges HTTP health → systemd watchdog
Type=notify
NotifyAccess=all
WatchdogSec=300
WorkingDirectory=/home/tca/rangebar-py
EnvironmentFile=/home/tca/rangebar-py/deploy/rangebar.env
EnvironmentFile=/home/tca/rangebar-py/deploy/sidecar.env
ExecStart=/home/tca/.cargo/bin/healthdog \
  --healthcheck-url http://localhost:8081/health \
  -- /home/tca/rangebar-py/.venv/bin/python3 /home/tca/rangebar-py/scripts/streaming_sidecar.py --from-env

# Restart policy (Issue #107 P3)
Restart=always
RestartSec=10

# OOM protection (Issue #107 P4 - prevents batch jobs from starving sidecar)
MemoryHigh=2G
MemoryMax=2.5G
MemorySwapMax=0
OOMScoreAdjust=-500

# Process limits
LimitNOFILE=65536
LimitNPROC=4096

# Standard output/error
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rangebar-sidecar

[Install]
WantedBy=default.target
