# deploy/rangebar-sidecar.service
# Issue #107: systemd user service for Layer 3 streaming sidecar.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp deploy/rangebar-sidecar.service ~/.config/systemd/user/
#   loginctl enable-linger tca
#   systemctl --user daemon-reload
#   systemctl --user enable --now rangebar-sidecar

[Unit]
Description=rangebar-py streaming sidecar (Layer 3)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tca/rangebar-py
ExecStart=/home/tca/rangebar-py/.venv/bin/python3 scripts/streaming_sidecar.py --from-env
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=5
MemoryMax=2G
MemorySwapMax=0
EnvironmentFile=/home/tca/rangebar-py/deploy/sidecar.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rangebar-sidecar

[Install]
WantedBy=default.target
